# bike_crash_bot.service
[Unit]
Description=Bike Crash Tracker Bot
After=network.target

[Service]
Type=simple
User=bike_bot
Group=bike_bot
WorkingDirectory=/opt/bike_crash_bot
Environment=BOT_TOKEN=1234567890:AAFxxxxxxxxxxxxxxxxxxxxx
Environment=DB_PATH=/var/lib/bike_crash_bot/state.db
ExecStart=/opt/bike_crash_bot/venv/bin/python /opt/bike_crash_bot/bike_crash_bot.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target